# puppet/manifests/files/uwsgi/uwsgi.service.erb
# REF: https://www.digitalocean.com/community/tutorials/how-to-serve-django-applications-with-uwsgi-and-nginx-on-centos-7

[Unit]
Description=uWSGI Emperor

[Service]
ExecStart=/usr/local/bin/uwsgi --die-on-term --emperor /etc/uwsgi/apps-enabled --logto <%= @logs_path %>/emperor.log
Restart=always
KillSignal=SIGQUIT
Type=notify
NotifyAccess=all

[Install]
WantedBy=multi-user.target
